# relay-public/clavisvault-relay.service
[Unit]
Description=ClavisVault Relay Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=clavisvault
Group=clavisvault
ExecStart=/usr/local/bin/clavisvault-relay --bind 0.0.0.0:51820
Restart=always
RestartSec=2
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
